# modules/uni-chop/systemd/unichop-servmon-int-bg2.service
[Unit]
Description = unichop monitoring server for INT
Requires = saftbus.service
After = saftbus.service

[Service]
Type = simple
Environment="DIM_DNS_NODE=asl105"
ExecStart = unichop-serv-mon tr0 -d int

[Install]
WantedBy = multi-user.target