# modules/b2b/systemd/b2b-analyzer-int-sis18.service
[Unit]
Description = b2b data analysis for INT SIS18

[Service]
Type = simple
Environment="DIM_DNS_NODE=asl105"
ExecStart = b2b-analyzer int_sis18

[Install]
WantedBy = multi-user.target